# pi/stage-openframe/01-server/files/openframe-server.service
[Unit]
Description=OpenFrame Docker Compose Stack
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service
ConditionPathExists=/home/pi/openframe/.env

[Service]
Type=oneshot
RemainAfterExit=yes
User=pi
Group=pi
WorkingDirectory=/home/pi/openframe
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
